# files/system/etc/skel/.config/autostart/redshift.desktop
[Desktop Entry]
Type=Application
Name=Redshift
Exec=/usr/bin/redshift -O 5000
